# user-systemd-services/joystickwake.service
[Unit]
Description=Joystickwake

[Service]
Type=simple
ExecStart=/usr/bin/joystickwake
Restart=on-failure

[Install]
WantedBy=default.target
