# user-systemd-services/test.service
[Service]
ExecStart=/bin/echo $RESTIC_PASSWORD

