# user-systemd-services/backup.service.d/override.conf
[Service]
Environment="NOTIFY_EMAIL_ADDRESS=weedleeedle@gmail.com"

